# .ralph/systemd/opencode-ralph.service
[Unit]
Description=Opencode Server for PGTuskMaster Ralph
Documentation=https://opencode.ai
After=network.target

[Service]
Type=simple
ExecStart=/home/joshazimullah.linux/.bun/bin/opencode serve
Restart=always
RestartSec=5
Environment="HOME=%h"
Environment="USER=%u"
WorkingDirectory=%h/work_mounts/projects/postgres_operator/PGTuskMaster/ElixirPGTuskMaster

[Install]
WantedBy=default.target
