# .ralph/systemd/ralph-event-watch.service
[Unit]
Description=Event Watcher for PGTuskMaster Ralph
Documentation=https://opencode.ai
After=network.target opencode-ralph.service
Requires=opencode-ralph.service

[Service]
Type=simple
ExecStart=/bin/bash /home/joshazimullah.linux/work_mounts/projects/postgres_operator/PGTuskMaster/ElixirPGTuskMaster/.ralph/event_watch.sh
Restart=always
RestartSec=5
Environment="HOME=%h"
Environment="USER=%u"
WorkingDirectory=%h/work_mounts/projects/postgres_operator/PGTuskMaster/ElixirPGTuskMaster

[Install]
WantedBy=default.target
